# build/microvm/inittab
::sysinit:/etc/rc.local
::respawn:/usr/local/bin/lb-agent
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r
