# images/mkosi.profiles/ecs/mkosi.extra/etc/systemd/system/mulga-cdi-generate.service.d/10-ecs-ordering.conf
# The CDI unit ships in gpu-nvidia without a Before=, because what it must
# precede differs per image. ECS needs the spec on disk before containerd
# starts, and containerd's unit comes from the Ubuntu package so it cannot
# carry an After= of its own. ecs-agent.service already orders itself after the
# CDI unit; it is repeated here so this drop-in states the full ordering rather
# than half of it split across two files.
#
# Inert on a CPU ECS image: with no gpu-nvidia profile there is no
# mulga-cdi-generate.service, and systemd ignores a drop-in for a unit that
# does not exist.
[Unit]
Before=containerd.service ecs-agent.service
