# images/mkosi.extra/etc/systemd/system/systemd-networkd-wait-online.service.d/override.conf
# wait-online blocks network-online.target until EVERY managed link is up, and
# it is not an opt-in: docker.service pulls the target in with Wants=, so any
# image carrying docker inherits the wait. A guest whose links are slow or
# absent therefore runs the unit to its full 120s default, and multi-user.target
# — the login prompt included — queues up behind it.
#
# --any returns as soon as one link is up, which is what a guest with a single
# NIC actually needs, and the short timeout bounds the case where none appears.
#
# Baked into the image rather than delivered as a first-boot drop-in: an image
# built from scratch has no cloud-init to write it.
[Service]
ExecStart=
ExecStart=/usr/lib/systemd/systemd-networkd-wait-online --any --timeout=15
