# build/systemd/spinifex-nats-watchdog.timer
[Unit]
Description=Spinifex NATS JetStream liveness watchdog timer
PartOf=spinifex.target

[Timer]
# Delay the first probe so NATS startup + JetStream cluster quorum settle (a
# forming cluster reports JS not-current transiently, which is not a fault).
OnBootSec=60s
OnUnitActiveSec=30s
AccuracySec=5s

[Install]
WantedBy=timers.target
